# thecafedisite/deploy/gunicorn.socket
[Unit]
Description=gunicorn socket for thecafedisco

[Socket]
ListenStream=/run/gunicorn.sock

[Install]
WantedBy=sockets.target
